# agent/launchd/com.cortex.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortex.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/dfonnegrag/.nvm/versions/node/v22.12.0/bin/node</string>
    <string>--env-file=/Users/dfonnegrag/Projects/cortex/.env</string>
    <string>--import=tsx</string>
    <string>/Users/dfonnegrag/Projects/cortex/agent/src/index.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/dfonnegrag/Projects/cortex</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/cortex-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cortex-daemon-error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>HOME</key>
    <string>/Users/dfonnegrag</string>
    <key>PATH</key>
    <string>/Users/dfonnegrag/.nvm/versions/node/v22.12.0/bin:/Users/dfonnegrag/Projects/cortex/node_modules/.bin:/Users/dfonnegrag/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>CORTEX_API_URL</key>
    <string>https://REPLACE_WITH_VERCEL_URL.vercel.app</string>
    <key>CORTEX_API_KEY</key>
    <string>REPLACE_WITH_API_KEY_FROM_VERCEL</string>
    <key>WATCH_PATHS</key>
    <string>/Users/dfonnegrag/Downloads,/Users/dfonnegrag/Documents</string>
  </dict>
</dict>
</plist>
